# systemd/phonesystem.service
[Unit]
Description=Raspberry Pi Phone System - IFB/PL
After=network.target sound.target pulseaudio.service
Wants=pulseaudio.service

[Service]
Type=simple
User=procomm
Group=procomm
WorkingDirectory=/home/procomm/ProComm
Environment="QT_QPA_PLATFORM=linuxfb"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/procomm/.Xauthority"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
ExecStart=/usr/bin/python3 /home/procomm/ProComm/main.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/procomm/ProComm/logs/systemd.log
StandardError=append:/home/procomm/ProComm/logs/systemd.log

# Resource limits
LimitRTPRIO=95
LimitMEMLOCK=infinity
Nice=-10

[Install]
WantedBy=graphical.target
