# systemd/phonesystem.service
[Unit]
Description=Raspberry Pi Phone System - IFB/PL
After=network.target sound.target
DefaultDependencies=no

[Service]
Type=simple
User=procomm
Group=procomm
WorkingDirectory=/home/procomm/ProComm
Environment="QT_QPA_PLATFORM=linuxfb"
Environment="QT_QPA_FB_DRM=1"
Environment="QT_QPA_EGLFS_ALWAYS_SET_MODE=1"
ExecStart=/usr/bin/python3 /home/procomm/ProComm/main.py
Restart=always
RestartSec=5
StandardOutput=append:/home/procomm/ProComm/logs/systemd.log
StandardError=append:/home/procomm/ProComm/logs/systemd.log

# Resource limits
LimitRTPRIO=95
LimitMEMLOCK=infinity
Nice=-10

# Ensure framebuffer access
SupplementaryGroups=video audio

[Install]
WantedBy=multi-user.target
